# services/gpu/deploy/supervisord.conf
; LiveAvatarStream GPU plane process manager.
;
; Runs seven GPU services plus engine-three (Node) and the nginx gateway (0.0.0.0:8080).
; Ports: avatar-build 8011, image-gen 8002, voice 8003, avatar-video 8004,
; finishing 8005, realtime 8006, engine-three 8090.
; module imports (app, models, pipeline, ...) resolve; las_common is on
; PYTHONPATH (exported by start.sh). Logs go to stdout/stderr so RunPod captures
; them. autorestart keeps a crashed worker (e.g. OOM) coming back.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:avatar-build]
command=uvicorn app:app --host 127.0.0.1 --port 8011
directory=%(ENV_GPU_DIR)s/avatar-build
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:image-gen]
command=uvicorn app:app --host 127.0.0.1 --port 8002
directory=%(ENV_GPU_DIR)s/image-gen
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:voice]
command=uvicorn app:app --host 127.0.0.1 --port 8003
directory=%(ENV_GPU_DIR)s/voice
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:avatar-video]
command=uvicorn app:app --host 127.0.0.1 --port 8004
directory=%(ENV_GPU_DIR)s/avatar-video
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:finishing]
command=uvicorn app:app --host 127.0.0.1 --port 8005
directory=%(ENV_GPU_DIR)s/finishing
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:realtime]
command=uvicorn app:app --host 127.0.0.1 --port 8006
directory=%(ENV_GPU_DIR)s/realtime
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

[program:engine-three]
command=node dist/index.js
directory=%(ENV_ENGINE_THREE_DIR)s
environment=DISPLAY=":99",PORT="8090",LAS_ASSETS_DIR="%(ENV_ENGINE_THREE_DIR)s/assets",CONTROL_API_URL="%(ENV_CONTROL_API_URL)s",INTERNAL_SERVICE_TOKEN="%(ENV_INTERNAL_SERVICE_TOKEN)s",R2_ENDPOINT="%(ENV_R2_ENDPOINT)s",R2_ACCESS_KEY_ID="%(ENV_R2_ACCESS_KEY_ID)s",R2_SECRET_ACCESS_KEY="%(ENV_R2_SECRET_ACCESS_KEY)s",R2_OUTPUTS_BUCKET="%(ENV_R2_OUTPUTS_BUCKET)s",RENDER_PROFILE="%(ENV_RENDER_PROFILE)s",LIPSYNC_MODE="%(ENV_LIPSYNC_MODE)s",MONTAGE_MODE="%(ENV_MONTAGE_MODE)s"
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15

[program:nginx]
command=nginx -c %(ENV_GPU_DIR)s/deploy/nginx.conf
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Start the gateway after the workers have had a moment to bind.
priority=50
