# services/gpu/deploy/supervisord.conf
; LiveAvatarStream GPU plane process manager.
;
; Runs the six FastAPI services (uvicorn) plus the nginx gateway (0.0.0.0:8080).
; Ports: avatar-build 8011, image-gen 8002, voice 8003, avatar-video 8004,
; finishing 8005, realtime 8006 (avatar-build moved off 8001 to avoid a pod port
; clash — see deploy/POD_SETUP.md §5). Each service runs from its own directory so local
; module imports (app, models, pipeline, ...) resolve; las_common is on
; PYTHONPATH (exported by start.sh). Logs go to stdout/stderr so RunPod captures
; them. autorestart keeps a crashed worker (e.g. OOM) coming back.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:avatar-build]
command=uvicorn app:app --host 127.0.0.1 --port 8011
directory=%(ENV_GPU_DIR)s/avatar-build
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:image-gen]
command=uvicorn app:app --host 127.0.0.1 --port 8002
directory=%(ENV_GPU_DIR)s/image-gen
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:voice]
command=uvicorn app:app --host 127.0.0.1 --port 8003
directory=%(ENV_GPU_DIR)s/voice
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:avatar-video]
command=uvicorn app:app --host 127.0.0.1 --port 8004
directory=%(ENV_GPU_DIR)s/avatar-video
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:finishing]
command=uvicorn app:app --host 127.0.0.1 --port 8005
directory=%(ENV_GPU_DIR)s/finishing
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:realtime]
command=uvicorn app:app --host 127.0.0.1 --port 8006
directory=%(ENV_GPU_DIR)s/realtime
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -c %(ENV_GPU_DIR)s/deploy/nginx.conf
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Start the gateway after the workers have had a moment to bind.
priority=50
